# dotpi_root/dist/@dotpi/avahi-monitor/dotpi-avahi-monitor.service
[Unit]
Description=avahi monitor
After=avahi-daemon.service

[Service]
Type=exec
ExecStart=/opt/dotpi/lib/dotpi_modules/@dotpi/avahi-monitor-module/dotpi-avahi-monitor-command
Restart=always

# wait before restarting
RestartSec=10s

# Do not stop to try restarting service
StartLimitBurst=0

[Install]
WantedBy=multi-user.target

